# supervisord.conf
[supervisord]
nodaemon=true

[program:report-generator]
command=/bin/bash -c "/usr/local/bin/report-generator"
